# zabbix/supervisord_zabbix.conf
[program:zabbix]
command=docker-compose -f /home/ubuntu/zabbix/zabbix.yaml up
